# scripts/com.teamworkmissive.connector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique identifier for this launch agent -->
    <key>Label</key>
    <string>com.teamworkmissive.connector</string>
    
    <!-- Working directory (replace with your actual project path) -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/TeamworkMissiveConnector</string>
    
    <!-- Command to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>scripts/run_local.sh</string>
    </array>
    
    <!-- Run at load (login) -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Keep the process alive (restart if it crashes) -->
    <key>KeepAlive</key>
    <true/>
    
    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/TeamworkMissiveConnector/logs/launchd.out.log</string>
    
    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/TeamworkMissiveConnector/logs/launchd.err.log</string>
    
    <!-- Environment variables (if needed) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <!-- Nice value (lower priority = 10, normal = 0) -->
    <key>Nice</key>
    <integer>0</integer>
</dict>
</plist>

